# deployment/supervisor/chirp-tasks.conf
[program:chirp-tasks]
command=/home/ubuntu/chirp/venv/bin/rq worker chirp-tasks
numprocs=1
directory=/home/ubuntu/chirp
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
